# Advanced/Server/central-sensor-dashboard.service
[Unit]
Description=Central Sensor Dashboard Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=UsernameHere
WorkingDirectory=/Path/To/file/
ExecStart=/usr/bin/python3 /Path/To/file/sensor_api.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional)
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target